# opie/resources/com.opie.relay.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  LaunchAgent for the Opie voice relay. Rendered from this template by
  opie/service.py (the GUI's "Autostart at login" toggle) — paths below are filled
  in with the venv Python, your config file, and a log file.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>{{LABEL}}</string>

    <key>ProgramArguments</key>
    <array>
        <string>{{PYTHON}}</string>
        <string>-m</string>
        <string>opie</string>
        <string>--config</string>
        <string>{{CONFIG}}</string>
    </array>

    <key>WorkingDirectory</key>
    <string>{{WORKDIR}}</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>{{STDOUT}}</string>
    <key>StandardErrorPath</key>
    <string>{{STDERR}}</string>
</dict>
</plist>
